# systemd/eda-server.service
[Unit]
Description=Event-Driven Ansible — DDIL Identity Rulebook
After=identity-service.service skupper-init.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ansible-rulebook \
    --rulebook /opt/eda/rulebook.yml \
    --inventory localhost, \
    --verbose
Environment=RHSI_HEALTH_URL=http://localhost:9090/healthz
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
